# meta/mpu6050.service
[Service]
Type=simple
RemainAfterExit=no
Restart=always
ExecStart=mpu6050

[Install]
WantedBy=multi-user.target
